# home/dot/config/systemd/user/paru-cache-clean.service
[Unit]
Description=Clean paru cache

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c 'rm -rf "$HOME/.cache/paru/"*'
